# templates/ovs-restart.j2
#!/bin/sh
#
# chkconfig: - 99 01
#
case "$1" in
  start)

    sleep 20
    systemctl restart openvswitch

    ;;
esac

quit 0
